# docs-hub.service
[Unit]
Description=Docs Hub
After=network.target

[Service]
Type=simple
User=docshub
Group=docshub
WorkingDirectory=/opt/docs-hub
Environment=ADDR=127.0.0.1:8080
Environment=DATA_FILE=/var/lib/docs-hub/storage.json
Environment=ADMIN_USER=admin
Environment=ADMIN_PASSWORD=admin123
Environment=DOCS_HUB_MODE=server
ExecStart=/opt/docs-hub/docs-hub serve
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/docs-hub

[Install]
WantedBy=multi-user.target
